# worker-node/root.crontab
@reboot /root/bin/phone_home > /var/log/phone_home.log
